# services/scripts/coinbreakr.service
[Unit]
Description=CoinBreakr Docker Compose Service
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/coinbreakr
EnvironmentFile=/opt/coinbreakr/.env
ExecStartPre=/bin/bash -c 'until docker info >/dev/null 2>&1; do echo "Waiting for Docker daemon..."; sleep 2; done'
ExecStart=/usr/local/bin/docker-compose -f /opt/coinbreakr/docker-compose.yml up --build
ExecStop=/usr/local/bin/docker-compose -f /opt/coinbreakr/docker-compose.yml down
Restart=always
RestartSec=10s
TimeoutStartSec=600
TimeoutStopSec=120
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
